# shadowsocks/rss.free/systemd-shadowsocks-rowan-chk.service
[Unit]
Description=free shadow socks service

[Service]
Type=simple
ExecStart=/root/bin/log_time.sh
