# shadowsocks/heroku.websocks/systemd-rss-heroku.service
[Unit]
Description=heroku free shadow socks service

[Service]
Type=forking
ExecStart=/root/bin/log_time.sh

[Install]
WantedBy=multi-user.target
