# shadowsocks/systemd-rss-heroku-chk.service
[Unit]
Description=rss-heroku-chk service

[Service]
Type=simple
ExecStart=/root/bin/log_time.sh
